# sh/rails_serverd
#!/bin/bash
#chkconfig: - 85 15
#description: nginx is a World Wide Web server. It is used to serve
. /home/qishanqing/myscript/sh/server

source  Rails
