# jenkins/jenkinsd
#!/bin/bash
#chkconfig: - 85 15
#description: nginx is a World Wide Web server. It is used to serve

su - qishanqing <<EOF
cd /home/qishanqing/myscript/jenkins
java -jar jenkins.war --httpPort=8080&
EOF
